# dev/com.aiki.otel-receive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aiki.otel-receive</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/aiki</string>
        <string>otel-receive</string>
    </array>

    <!-- Socket activation: pass incoming connection as stdin -->
    <key>Sockets</key>
    <dict>
        <key>Listeners</key>
        <dict>
            <key>SockServiceName</key>
            <string>19876</string>
            <key>SockNodeName</key>
            <string>127.0.0.1</string>
            <key>SockType</key>
            <string>stream</string>
        </dict>
    </dict>

    <!-- inetd-style: stdin/stdout are the socket -->
    <key>inetdCompatibility</key>
    <dict>
        <key>Wait</key>
        <false/>
    </dict>

    <!-- Logging -->
    <key>StandardErrorPath</key>
    <string>/tmp/aiki-otel-receive.err</string>

    <!-- Process spawning settings -->
    <key>SessionCreate</key>
    <false/>

    <!-- Don't keep running - only launch on socket activation -->
    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
